# cron/crontab
#crontab
#0-59/1  *  *  *  *		root	*alerts/test_mail_piry
00   18  *  *  *		root	*alerts/alertUsersLastRecommendations
