# cron/crontab
#crontab
#0-59/1  *  *  *  *		root	*alerts/test_mail_piry
00   11  *  *  *		root	*alerts/alertUsersLastRecommendations
