# cron/crontab
#crontab
#0-59/1  *  *  *  *		root	*alerts/test_mail_piry
30   09  *  *  *		root	*alerts/alertUsersLastRecommendations
