# cron/crontab
#crontab
00   16  *  *  *		root	*alerts/alertUsersLastRecommendations
