# utils/services/check-run-mailqueues.timer
[Unit]
Description=Run check run-mailqueues service each minute.

[Timer]
OnBootSec=5min
OnCalendar=*:0/1

[Install]
WantedBy=timers.target